# pkg/smokegios.cron
#
# Check for configuration changes in Nagios to apply to Smokeping
#

01 * * * *     root test -x /usr/bin/smokegios -c /etc/smokeping/smokegios.conf

